# src/Resources/ubuntu/supervisor/app.conf
[program:unit3d-queue]
process_name=%(program_name)s_%(process_num)02d
command=php {{INSTALLDIR}}/artisan queue:work --sleep=3 --tries=3
startsecs = 0
autostart=true
autorestart=true
user={{WEBUSER}}
numprocs=2

[program:unit3d-chat-server]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/node /usr/bin/laravel-echo-server start --dir={{INSTALLDIR}}
autostart=true
autorestart=true
user={{WEBUSER}}
numprocs=1
